# install/files/systemd/node-aleo-snarkos.service.template
[Unit]
Description=ALEO snarkOS Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${POOL_USER}
Group=${POOL_USER}

# Ensure log and data directories exist before starting
ExecStartPre=/bin/mkdir -p ${ALEO_DIR}/logs
ExecStartPre=/bin/mkdir -p ${ALEO_DIR}/data
ExecStart=${ALEO_DIR}/bin/start-node.sh

Restart=on-failure
RestartSec=30
TimeoutStartSec=infinity
TimeoutStopSec=120

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
