# install/files/systemd/solo-pool-payments.service.template
[Unit]
Description=Solo Pool Payment Processor
After=network-online.target
Wants=network-online.target
# Payments can start early - it will gracefully handle unavailable wallets

[Service]
Type=simple
User=${POOL_USER}
Group=${POOL_USER}

WorkingDirectory=${PAYMENTS_DIR}
ExecStart=${PAYMENTS_DIR}/bin/solo-pool-payments ${PAYMENTS_DIR}/config/config.toml

Environment=RUST_LOG=info

Restart=on-failure
RestartSec=30

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
