# install/files/systemd/wallet-xtm.service.template
[Unit]
Description=Tari Wallet (Pool Payments)
After=network-online.target node-xtm-minotari.service
Wants=network-online.target
Requires=node-xtm-minotari.service

[Service]
Type=simple
User=${POOL_USER}
Group=${POOL_USER}

ExecStart=${TARI_DIR}/bin/start-wallet.sh
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=60

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
