# install-scripts/files/systemd/node-aleo-snarkos.service
[Unit]
Description=ALEO snarkOS Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${POOL_USER}
Group=${POOL_USER}

ExecStart=${ALEO_DIR}/start-node.sh

Restart=on-failure
RestartSec=30
TimeoutStartSec=infinity
TimeoutStopSec=120

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
