# install-scripts/files/systemd/node-xtm-minotari.service
[Unit]
Description=Minotari Node (Tari)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${POOL_USER}
Group=${POOL_USER}

ExecStart=${TARI_DIR}/bin/minotari_node --config=${TARI_DIR}/config/config.toml --base-path=${TARI_DIR}/data
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=on-failure
RestartSec=30
TimeoutStartSec=infinity
TimeoutStopSec=120

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
